# systemd/cadence-voice-daemon.service
[Unit]
Description=Cadence Voice Daemon — mic capture + Whisper + intent routing
After=openclaw-gateway.service pulseaudio.service
Wants=openclaw-gateway.service

[Service]
Type=simple
ExecStart=/home/linuxbrew/.linuxbrew/bin/python3 /home/rollan/.openclaw/workspace/jarvis-v5/runtime/voice/cadence_daemon.py --execute --device RDPSource --duration 4
Restart=on-failure
RestartSec=10
TimeoutStopSec=15
# control-group kills all child processes (parecord) on stop/restart
KillMode=control-group

Environment=HOME=/home/rollan
Environment=TMPDIR=/tmp
Environment=PULSE_SERVER=unix:/mnt/wslg/PulseServer
Environment=CADENCE_WHISPER_MODEL=small.en
Environment=CADENCE_PASSIVE_CHUNK_SEC=4.0
Environment=CADENCE_COMMAND_WINDOW_SECONDS=8.0

WorkingDirectory=/home/rollan/.openclaw/workspace/jarvis-v5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=cadence-voice

[Install]
WantedBy=default.target
