# systemd/user/openclaw-review-poller.service
[Unit]
Description=OpenClaw Discord Review Poller
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/linuxbrew/.linuxbrew/bin/python3 /home/rollan/.openclaw/workspace/jarvis-v5/scripts/discord_review_poller.py --once --verbose
WorkingDirectory=/home/rollan/.openclaw/workspace/jarvis-v5
Environment=HOME=/home/rollan
Environment=TMPDIR=/tmp
Environment=PYTHONDONTWRITEBYTECODE=1

[Install]
WantedBy=default.target
