# systemd/user/openclaw-kitt-paper.service
[Unit]
Description=OpenClaw Kitt Paper Trading — one bounded cycle
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/rollan/.openclaw/workspace/jarvis-v5
ExecStart=/home/rollan/.openclaw/workspace/jarvis-v5/.venv/bin/python3 /home/rollan/.openclaw/workspace/jarvis-v5/workspace/quant_infra/kitt/run_kitt_cycle.py
EnvironmentFile=/home/rollan/.openclaw/secrets.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-kitt-paper
TimeoutStartSec=120
